# systemd/hangar-supervisor.service
[Unit]
Description=Hangar PTY supervisor daemon
Before=hangar.service

[Service]
Type=simple
ExecStart=%h/Documents/hangar/backend/target/release/hangar-supervisor
Environment=RUST_LOG=info
Restart=always
RestartSec=1
# Only kill supervisor process on stop, not its PTY children
# Children get reparented to PID 1 until supervisor restarts
KillMode=process

[Install]
WantedBy=default.target
