# systemd/hangar.service
[Unit]
Description=Hangar backend daemon
After=network.target

[Service]
Type=simple
User=george
WorkingDirectory=/home/george/Documents/hangar
ExecStart=/home/george/Documents/hangar/backend/target/release/hangard
Environment=HANGAR_PORT=3000
Environment=RUST_LOG=info
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
