# systemd/ttyd-issue12.service
[Unit]
Description=ttyd for tmux session: issue12
After=network.target
Requires=network.target

[Service]
Type=simple
User=george
Group=george
Environment=TERM=xterm-256color
ExecStart=/usr/bin/ttyd -p 7684 -W -t titleFixed=issue12 -t fontSize=13 tmux attach -t issue12
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
